# systemd-user/.config/systemd/user/emacs.service.d/log.conf
[Service]
StandardOutput=append:%h/.emacs-stderr.log
StandardError=append:%h/.emacs-stderr.log
